# deploy/autodosie-web.service
[Unit]
Description=Autodosie web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autobot
Group=autobot
WorkingDirectory=/home/autobot/apps/autodosie_bot
EnvironmentFile=/home/autobot/apps/shared/autodosie_bot.env
ExecStart=/home/autobot/.venvs/autodosie_bot/bin/autodosie-web
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
